# SpeedDevice/rc.local
#!/bin/sh -e

logger RUN SPEED TEST ON BOOT

#execute firstboot.sh only once
if [ ! -e /root/firstboot_done ]; then
    if [ -e /root/firstboot.sh ]; then
        /root/firstboot.sh
    fi
    touch /root/firstboot_done
fi
exit 0

